# deploy/systemd/yolo11-rtsp.service.template
[Unit]
Description=YOLO11 RTSP Web Detector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
Environment=HOST=0.0.0.0
Environment=PORT=8000
EnvironmentFile=__ENV_FILE__
ExecStart=__UVICORN_BIN__ app.main:app --host ${HOST} --port ${PORT}
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
